# etc/supervisord.d/webserver.ini
[program:webserver]
command = python -m SimpleHTTPServer 80
directory = /opt/graphgist/
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
